# qbittorrent/s6-overlay/etc/services.d/log/finish
#!/usr/bin/with-contenv bash

kill -n 9 $(ps | grep "tail -F /data/data/logs/qbittorrent.log" | grep -v "grep" | awk '{print $1}')
