# qbittorrent/s6-overlay/etc/services.d/log/run
#!/usr/bin/with-contenv bash

sleep 5
cd /data/data/logs/
exec s6-setuidgid qbittorrent tail -F qbittorrent.log
