# chinesesubfinder/s6-overlay/etc/services.d/chinesesubfinder/run
#!/usr/bin/with-contenv sh

umask ${UMASK:-022}
cd /app
exec s6-setuidgid csf ./chinesesubfinder
