# subfinder/s6-overlay/etc/services.d/subfinder/run
#!/usr/bin/with-contenv sh

umask ${UMASK:-022}

exec s6-setuidgid subfinder loop -e ${INTERVAL} -- subfinder -c /config/subfinder.json /media
