# test/s6-overlay/etc/services.d/qbittorrent/run
#!/usr/bin/with-contenv bash

UMASK_SET=${UMASK_SET:-000}
umask $UMASK_SET

exec s6-setuidgid qbittorrent qbittorrent-nox --profile=${QBITTORRENT_HOME}
