# qbittorrent/s6-overlay/etc/services.d/crond/finish
#!/usr/bin/with-contenv bash

kill -n 9 $(ps | grep "crond -f" | grep -v "grep" | awk '{print $1}')
