# qbittorrent/s6-overlay/etc/services.d/qbittorrent/finish
#!/usr/bin/with-contenv bash

. /usr/local/bin/share
$cmd_curl ${api_url_base}/app/shutdown
